# rules/rules.d/99-bluetooth-headset-profile.rules
ACTION=="add", SUBSYSTEM=="input", ENV{PRODUCT}=="5/54c/cd3/452", RUN+="/usr/bin/su untermensch -c /home/untermensch/.local/bin/WH-1000XM3-notify"
ACTION=="add", SUBSYSTEM=="input", ENV{PRODUCT}=="5/0/0/0", RUN+="/usr/bin/su untermensch -c /home/untermensch/.local/bin/FreeBuds3-notify"
